# setup/worker/packages/nodeexporter-1.0/DEBIAN/postinst
#!/usr/bin/env bash

systemctl daemon-reload
systemctl enable node_exporter
systemctl start node_exporter
