# setup/worker/psexporter-1.0/etc/systemd/system/ps_exporter.service
[Unit]
Description=Ps Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/ps_exporter/ps_exporter.sh

[Install]
WantedBy=multi-user.target
