# setup/worker/psexporter-1.0/DEBIAN/postinst
#!/usr/bin/env bash

systemctl daemon-reload
systemctl enable ps_exporter
systemctl start ps_exporter
