# setup/worker/ipmiexporter-1.0/DEBIAN/postinst
#!/usr/bin/env bash

systemctl daemon-reload
systemctl enable ipmi_exporter
systemctl start ipmi_exporter
